# ops/systemd/tamagotchi-queue.service
[Unit]
Description=Tamagotchi Laravel Queue Worker
After=network.target redis.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/tamagotchi-api
ExecStart=/usr/bin/php artisan queue:work --sleep=3 --tries=3 --max-time=3600
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Environment
Environment=APP_ENV=production
Environment=HOME=/var/www

# Resource limits
LimitNOFILE=65536
MemoryLimit=512M

[Install]
WantedBy=multi-user.target
