# ops/systemd/tamagotchi-autodeploy.service
[Unit]
Description=Tamagotchi Auto Deploy Service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/tamagotchi-api
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStart=/bin/bash -lc '/opt/tamagotchi-api/ops/scripts/auto-deploy.sh'
Nice=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
